# docker/supervisord.conf
[supervisord]
nodaemon=true
user=root
loglevel=debug

[program:chargepal_bundle]
command=roslaunch chargepal_bundle startup.launch
autostart=true
autorestart=true
user=root
# stdout_logfile=/dev/fd/2
# stdout_logfile_maxbytes=0
# redirect_stderr=true

# [program:chargepal_behaviour_tree]
# command=rosrun chargepal_behaviour_tree chargepal_behaviour_tree
# autostart=true
# autorestart=true
# stdout_logfile=/dev/fd/1
# stdout_logfile_maxbytes=0
# redirect_stderr=true